# deploy/ansible/templates/smartcopy-agent.service.j2
[Unit]
Description=SmartCopy Agent - High-Performance File Transfer Agent
Documentation=https://github.com/smartcopy/smartcopy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ smartcopy_user }}
Group={{ smartcopy_group }}
EnvironmentFile=-{{ smartcopy_config_dir }}/smartcopy.env
ExecStart={{ smartcopy_install_dir }}/smartcopy agent \
    --protocol {{ smartcopy_agent_protocol }} \
{% if smartcopy_agent_protocol == 'tcp' %}
    --port {{ smartcopy_agent_port }} \
    --bind {{ smartcopy_agent_bind }}
{% endif %}

Restart=always
RestartSec=5
StandardOutput=append:{{ smartcopy_log_dir }}/agent.log
StandardError=append:{{ smartcopy_log_dir }}/agent.error.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths={{ smartcopy_data_dir }} {{ smartcopy_log_dir }}

# Resource limits
LimitNOFILE=1048576
LimitNPROC=65535
LimitMEMLOCK=infinity

# Performance tuning
Nice=-10
IOSchedulingClass=realtime
IOSchedulingPriority=0

[Install]
WantedBy=multi-user.target
